# ubuntu1804x64/crontab
@reboot root /usr/bin/dvb_proxy_delayed > /var/log/dvb_proxy

